# controller/snell-ufw-manager.service
[Unit]
Description=Snell UFW Manager
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/snell-ufw-manager/controller
Environment=SNELL_DB=/opt/snell-ufw-manager/controller/data/snell_manager.db
Environment=SNELL_CONFIG=/opt/snell-ufw-manager/controller/config.yaml
ExecStart=/opt/snell-ufw-manager/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8899 --workers 1
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
